# python-ai-service/supervisord.conf
[supervisord]
logfile=%(here)s/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=%(here)s/logs/supervisord.pid
nodaemon=false
directory=%(here)s

[unix_http_server]
file=%(here)s/logs/supervisor.sock

[supervisorctl]
serverurl=unix://%(here)s/logs/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:fastapi]
command=python -m uvicorn app.main:app --host 0.0.0.0 --port 8001
directory=%(here)s
autostart=true
autorestart=true
startretries=3
startsecs=3
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=%(here)s/logs/fastapi.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=PYTHONPATH="%(here)s/.."

[program:celery-worker]
command=python -m celery -A app.services.celery_app:celery_app worker --loglevel=info --concurrency=2
directory=%(here)s
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=30
stopsignal=TERM
redirect_stderr=true
stdout_logfile=%(here)s/logs/celery-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=PYTHONPATH="%(here)s/.."
